# deploy/systemd/openassistd.service
[Unit]
Description=OpenAssist Daemon
After=network.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=__OPENASSIST_INSTALL_DIR__
EnvironmentFile=__OPENASSIST_ENV_FILE__
Environment=OPENASSIST_ENV_FILE=__OPENASSIST_ENV_FILE__
ExecStart=__OPENASSIST_NODE_BIN__ __OPENASSIST_INSTALL_DIR__/apps/openassistd/dist/index.js run --config __OPENASSIST_CONFIG_PATH__
Restart=always
RestartSec=5
TimeoutStopSec=30

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=__OPENASSIST_INSTALL_DIR__ __OPENASSIST_RW_CONFIG_DIR__ %h/.local/state/openassist
LockPersonality=true
RestrictSUIDSGID=true

[Install]
WantedBy=default.target
